# extra/memcached/config/supervisor.conf
[supervisord]
nodaemon=true

[program:memcached]
command=/usr/bin/memcached
user=memcached
autostart=true
autorestart=true